# flask.ini
[program:flask]
command=/Users/achilles_xushy/.pyenv/versions/anaconda3-4.2.0/bin/gunicorn -c gunicorn_config.py achilles_flask_first:app
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
numprocs_start=1
stopwaitsecs=12
user=achilles_xushy
directory=/Users/achilles_xushy/PycharmProjects/achilles_flask_first
stdout_logfile=/Users/achilles_xushy/PycharmProjects/config_and_log/supervisor/flask.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/Users/achilles_xushy/PycharmProjects/config_and_log/supervisor/flaskerr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
